# infra/pulse-backup.service
[Unit]
Description=pulse backup (snapshot de SQLite DB)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/opt/pulse/backup.sh
