# infra/pulse-backup.service
[Unit]
Description=pulse backup (pg_dump PostgreSQL)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/opt/pulse_bogota/backup.sh
